# deploy/split-vpn-webui.service
[Unit]
Description=Split VPN Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/data/split-vpn-webui/split-vpn-webui --systemd
WorkingDirectory=/data/split-vpn-webui
Restart=on-failure
RestartSec=5s
# CAP_NET_ADMIN is required for ipset, iptables, and interface operations.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
# Redirect stdout/stderr to the persistent log directory.
StandardOutput=append:/data/split-vpn-webui/logs/split-vpn-webui.log
StandardError=append:/data/split-vpn-webui/logs/split-vpn-webui.log

[Install]
WantedBy=multi-user.target
